# WaarpPackaging/src/packaging/DEBS/control/prerm
service waarp-r66server stop all
# chkconfig --del waarp-r66server
exit 0
